# vendor/supervisord.conf
[supervisord]
nodaemon=true

[program:convos]
directory=/convos
command=/convos/vendor/bin/carton exec hypnotoad script/convos
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
autorestart=true

[program:convos-backend]
directory=/convos
command=/convos/vendor/bin/carton exec script/convos backend
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
autorestart=true

[program:redis]
command=/usr/bin/redis-server
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
